# project-2/ansible/enhanced/roles/python_web_app/files/flask.service
[Unit]
Description=Flask service to run python app.

[Service]
ExecStart=/app/flask_run.sh

[Install]
WantedBy=multi-user.target

